# deploy/systemd/rocksolid-backup.service
[Unit]
Description=RockSolidLicense Backup Job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/rocksolidlicense/rocksolid.env
Environment=PROJECT_ROOT=/opt/rocksolidlicense
Environment=ENV_FILE=/etc/rocksolidlicense/rocksolid.env
Environment=BACKUP_DIR=/var/lib/rocksolid/backups
Environment=LABEL=scheduled
ExecStart=/opt/rocksolidlicense/deploy/linux/backup-rocksolid.sh
User=rocksolid
Group=rocksolid
NoNewPrivileges=true
StateDirectory=rocksolid
